# dev-api-handler/api-consumer.service
[Unit]
Description=API Server Consumer
After=network-online.target

[Service]
Type=simple
User=ajz27
Group=ajz27
WorkingDirectory=/home/ajz27/techtitans/dev-api-handler/
ExecStart=/usr/bin/php /home/ajz27/techtitans/dev-api-handler/apiServer.php
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target